# osquery.conf
{
  "options": {
    "logger_plugin": "filesystem",
    "logger_path": "%USERPROFILE%\\Observability\\logs\\osquery",
    "logger_rotate": "true",
    "logger_rotate_max_files": "5",
    "disable_logging": "false",
    "schedule_splay_percent": "10",
    "utc": "true"
  },
  "schedule": {
    "running_processes": {
      "query": "SELECT pid, name, path, cmdline, start_time, uid FROM processes ORDER BY start_time DESC;",
      "interval": 60,
      "description": "Snapshot of all running processes"
    },
    "listening_ports": {
      "query": "SELECT p.name, p.pid, lp.port, lp.protocol, lp.address FROM listening_ports lp LEFT JOIN processes p ON lp.pid = p.pid WHERE lp.address != '' ORDER BY lp.port;",
      "interval": 60,
      "description": "Ports currently listening for connections"
    },
    "active_connections": {
      "query": "SELECT pos.pid, p.name, pos.local_address, pos.local_port, pos.remote_address, pos.remote_port, pos.state FROM process_open_sockets pos LEFT JOIN processes p ON pos.pid = p.pid WHERE pos.remote_port != 0;",
      "interval": 30,
      "description": "Active outbound and inbound TCP connections"
    },
    "logged_in_users": {
      "query": "SELECT liu.user, liu.host, liu.time, liu.type FROM logged_in_users liu;",
      "interval": 60,
      "description": "Currently logged in users"
    },
    "startup_items": {
      "query": "SELECT name, path, args, type, source FROM startup_items;",
      "interval": 300,
      "description": "Programs configured to run at startup (persistence)"
    },
    "scheduled_tasks": {
      "query": "SELECT name, action, path, enabled, state, next_run_time FROM scheduled_tasks WHERE enabled = 1;",
      "interval": 300,
      "description": "Enabled Windows scheduled tasks"
    },
    "security_products": {
      "query": "SELECT name, type, state, signatures_up_to_date FROM windows_security_products;",
      "interval": 300,
      "description": "Installed AV and firewall products and their status"
    },
    "loaded_drivers": {
      "query": "SELECT name, path, signed FROM drivers WHERE signed = 0;",
      "interval": 300,
      "description": "Unsigned kernel drivers (unusual, investigate these)"
    },
    "user_accounts": {
      "query": "SELECT username, uid, gid, description, directory FROM users;",
      "interval": 600,
      "description": "Local user accounts"
    }
  }
}
